# supervisord.conf
[program:celeryd]
command = python /home/dotcloud/current/manage.py celeryd -l info -E --settings=/home/dotcloud/current/settings_production
